# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# Mirrors CMD from https://github.com/gramps-project/gramps-web-api/blob/master/Dockerfile
# Update when upstream CMD changes. Last verified: 2026-06-20
[program:grampsweb]
command=gunicorn -w %(ENV_GUNICORN_NUM_WORKERS)s -b 0.0.0.0:5000 gramps_webapi.wsgi:app --timeout %(ENV_GUNICORN_TIMEOUT)s --limit-request-line 8190
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=GRAMPSWEB_CELERY_CONFIG__broker_url="redis://127.0.0.1:6379/0",GRAMPSWEB_CELERY_CONFIG__result_backend="redis://127.0.0.1:6379/0",GRAMPSWEB_RATELIMIT_STORAGE_URI="redis://127.0.0.1:6379/1"

[program:celery]
command=celery -A gramps_webapi.celery worker --loglevel=INFO --concurrency=2
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
environment=GRAMPSWEB_CELERY_CONFIG__broker_url="redis://127.0.0.1:6379/0",GRAMPSWEB_CELERY_CONFIG__result_backend="redis://127.0.0.1:6379/0",GRAMPSWEB_RATELIMIT_STORAGE_URI="redis://127.0.0.1:6379/1"
